# supervisor/sgd-workers.conf
# supervisor/sgd-workers.conf — Alternativa a contenedores separados
# Si prefieres un solo contenedor con supervisor en vez de 3 contenedores worker,
# copia este archivo a php/ y úsalo en el Dockerfile.
# Por defecto usamos contenedores separados (más simple de monitorear en Portainer).

[program:worker-default]
command=php artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=86400
directory=/var/www/html
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/www/html/storage/logs/worker-default.log
stdout_logfile_maxbytes=10MB

[program:worker-pqrsd]
command=php artisan queue:work --queue=pqrsd-ai --sleep=3 --tries=3 --max-time=86400
directory=/var/www/html
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/www/html/storage/logs/worker-pqrsd.log
stdout_logfile_maxbytes=10MB

[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/html
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/www/html/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
